# roles/ct_runtime_vllm/templates/vllm.service.j2
[Unit]
Description=vLLM OpenAI-compatible API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ct_runtime_vllm_user }}
Group={{ ct_runtime_vllm_group }}
WorkingDirectory={{ ct_runtime_vllm_workdir }}
EnvironmentFile={{ ct_runtime_vllm_env_file }}
ExecStart={{ ct_runtime_vllm_venv }}/bin/vllm serve ${VLLM_MODEL} --host ${VLLM_BIND_HOST} --port ${VLLM_PORT} --device ${VLLM_DEVICE} --tensor-parallel-size ${VLLM_TENSOR_PARALLEL_SIZE} --pipeline-parallel-size ${VLLM_PIPELINE_PARALLEL_SIZE} --gpu-memory-utilization ${VLLM_GPU_MEMORY_UTILIZATION} --dtype ${VLLM_DTYPE} --kv-cache-dtype ${VLLM_KV_CACHE_DTYPE} --max-model-len ${VLLM_MAX_MODEL_LEN} --max-num-seqs ${VLLM_MAX_NUM_SEQS} --max-num-batched-tokens ${VLLM_MAX_NUM_BATCHED_TOKENS}{% if ct_runtime_vllm_tool_call_parser | length > 0 %} --tool-call-parser ${VLLM_TOOL_CALL_PARSER}{% endif %}{% if ct_runtime_vllm_reasoning_parser | length > 0 %} --reasoning-parser ${VLLM_REASONING_PARSER}{% endif %}{% if ct_runtime_vllm_served_model_name | length > 0 %} --served-model-name ${VLLM_SERVED_MODEL_NAME}{% endif %} ${VLLM_EXTRA_ARGS}
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
