# roles/ct_runtime_vllm/templates/vllm.service.j2
[Unit]
Description=vLLM OpenAI-compatible API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ct_runtime_vllm_user }}
Group={{ ct_runtime_vllm_group }}
WorkingDirectory={{ ct_runtime_vllm_workdir }}
EnvironmentFile={{ ct_runtime_vllm_env_file }}
ExecStart={{ ct_runtime_vllm_venv }}/bin/vllm serve ${INFER_VLLM_MODEL} --host ${INFER_VLLM_BIND_HOST} --port ${INFER_VLLM_PORT}{% if ct_runtime_vllm_device != 'cuda' or ct_runtime_vllm_force_device_arg %} --device ${INFER_VLLM_DEVICE}{% endif %} --tensor-parallel-size ${INFER_VLLM_TENSOR_PARALLEL_SIZE} --pipeline-parallel-size ${INFER_VLLM_PIPELINE_PARALLEL_SIZE} --gpu-memory-utilization ${INFER_VLLM_GPU_MEMORY_UTILIZATION} --dtype ${INFER_VLLM_DTYPE} --kv-cache-dtype ${INFER_VLLM_KV_CACHE_DTYPE} --max-model-len ${INFER_VLLM_MAX_MODEL_LEN} --max-num-seqs ${INFER_VLLM_MAX_NUM_SEQS} --max-num-batched-tokens ${INFER_VLLM_MAX_NUM_BATCHED_TOKENS}{% if ct_runtime_vllm_tool_call_parser | length > 0 %} --tool-call-parser ${INFER_VLLM_TOOL_CALL_PARSER}{% endif %}{% if ct_runtime_vllm_reasoning_parser | length > 0 %} --reasoning-parser ${INFER_VLLM_REASONING_PARSER}{% endif %}{% if ct_runtime_vllm_pass_served_model_name and (ct_runtime_vllm_served_model_name_effective | default('') | length > 0) %} --served-model-name ${INFER_VLLM_SERVED_MODEL_NAME}{% endif %}{% if ct_runtime_vllm_extra_args | length > 0 %} ${INFER_VLLM_EXTRA_ARGS}{% endif %}
{% if ct_runtime_vllm_systemd_log_to_files %}
StandardOutput=append:{{ ct_runtime_vllm_stdout_log }}
StandardError=append:{{ ct_runtime_vllm_stderr_log }}
{% endif %}
Restart={{ ct_runtime_vllm_systemd_restart }}
RestartSec={{ ct_runtime_vllm_systemd_restart_sec }}
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
