# roles/ct_runtime_vllm/templates/vllm.service.j2
[Unit]
Description=vLLM OpenAI-compatible API server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ ct_runtime_vllm_user }}
Group={{ ct_runtime_vllm_group }}
WorkingDirectory={{ ct_runtime_vllm_workdir }}
EnvironmentFile={{ ct_runtime_vllm_env_file }}
ExecStart={{ ct_runtime_vllm_venv }}/bin/vllm serve ${VLLM_MODEL} --host ${VLLM_BIND_HOST} --port ${VLLM_PORT} --tensor-parallel-size ${VLLM_TENSOR_PARALLEL_SIZE} --max-model-len ${VLLM_MAX_MODEL_LEN} --gpu-memory-utilization ${VLLM_GPU_MEMORY_UTILIZATION} ${VLLM_EXTRA_ARGS}
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
